# amazon-s3-backup/rootfs/etc/services.d/amazon-s3-backup/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Amazon S3 Backup
# Runs amazon-s3-backup.py script
# ==============================================================================

bashio::log.info "Starting Amazon S3 Backup..."

AWS_ACCESS_KEY_ID="$(bashio::config 'aws_access_key')"
export AWS_ACCESS_KEY_ID

AWS_SECRET_ACCESS_KEY="$(bashio::config 'aws_secret_access_key')"
export AWS_SECRET_ACCESS_KEY

bucket_name="$(bashio::config 'bucket_name')"
export bucket_name

bucket_region="$(bashio::config 'bucket_region')"
export bucket_region

storage_class="$(bashio::config 'storage_class')"
export storage_class

upload_missing_files="$(bashio::config 'upload_missing_files')"
export upload_missing_files

keep_local_snapshots="$(bashio::config 'keep_local_snapshots')"
export keep_local_snapshots

export monitor_path="/backup"

echo "Я работаю"
echo "$AWS_ACCESS_KEY_ID"
echo "$AWS_SECRET_ACCESS_KEY"

exec python3 -u /usr/bin/amazon-s3-backup/amazon-s3-backup.py >&2